# bwapp_front.conf
[supervisord]
nodaemon=true

[program:apache2]
command=/etc/init.d/apache2 restart
